# ops/systemd/iot-sim-ops-api.service
[Unit]
Description=IoT SIM OPS API (uvicorn)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=fenrir
WorkingDirectory=/home/fenrir/iot-sim-ops/api/mock-fastapi
# 可选：.env，不存在也不报错（前面有 - 号）
EnvironmentFile=-/home/fenrir/iot-sim-ops/api/mock-fastapi/.env
# 用你项目里 venv 的 uvicorn 启动
ExecStart=/home/fenrir/iot-sim-ops/api/mock-fastapi/.venv/bin/uvicorn app:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=3
# 补 PATH，避免某些库找不到可执行文件
Environment=PATH=/home/fenrir/iot-sim-ops/api/mock-fastapi/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin

[Install]
WantedBy=multi-user.target
