# ops/systemd/iot-sim-ops-api.service
# 做成用户级 systemd 服务（长期后台跑）

mkdir -p ~/.config/systemd/user
cat > ~/.config/systemd/user/iot-sim-ops-api.service <<'EOF'
[Unit]
Description=iot-sim-ops FastAPI (user)

[Service]
WorkingDirectory=/home/fenrir/iot-sim-ops-api
EnvironmentFile=/home/fenrir/iot-sim-ops-api/.env
ExecStart=/home/fenrir/iot-sim-ops-api/.venv/bin/uvicorn app:app --host 0.0.0.0 --port 8000
Restart=on-failure

[Install]
WantedBy=default.target
EOF

systemctl --user daemon-reload
systemctl --user enable --now iot-sim-ops-api.service
systemctl --user status iot-sim-ops-api.service
